# plugins/lando-proxy/proxy/root/etc/services.d/hipache/run
#!/usr/bin/with-contenv sh

exec s6-setuidgid root /usr/local/bin/hipache -c /etc/hipache/config.json;
