# plugins/lando-proxy/proxy/root/etc/services.d/0-keys/run
#!/usr/bin/with-contenv sh

exec s6-setuidgid root /config/keys.sh;
