# plugins/lando-proxy/proxy/root/etc/services.d/redis/run
#!/usr/bin/with-contenv sh

exec s6-setuidgid root redis-server /etc/redis/redis.conf;
